# examples/docker/pyircbot.conf
[program:pyircbot3]
command=pyircbot -c config.json
autorestart=true
user=bot
directory=/srv/bot/
